# dot_config/systemd/user/noctalia.service
[Unit]
Description=Noctalia shell
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Environment=WAYLAND_DISPLAY=wayland-1
ExecStartPre=/bin/bash -c 'until [ -S "${XDG_RUNTIME_DIR}/wayland-1" ]; do sleep 0.5; done'
ExecStart=noctalia
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical-session.target
